# willys-shopping-list/rootfs/etc/s6-overlay/s6-rc.d/willys-sync/run
#!/usr/bin/with-contenv bash
set -e

echo "[willys] Startar..."

CONFIG_PATH=/data/options.json

if [ ! -f "$CONFIG_PATH" ]; then
    echo "[willys] FATAL: $CONFIG_PATH saknas"
    exit 1
fi

export WILLYS_USERNAME="$(jq -r '.willys_username' $CONFIG_PATH)"
export WILLYS_PASSWORD="$(jq -r '.willys_password' $CONFIG_PATH)"
export WILLYS_LIST_ID="$(jq -r '.willys_list_id // empty' $CONFIG_PATH)"
export SYNC_INTERVAL="$(jq -r '.sync_interval_seconds // 30' $CONFIG_PATH)"
export LOG_LEVEL="$(jq -r '.log_level // "info"' $CONFIG_PATH)"
export COOKIE_DIR="/share/willys"
export HA_URL="http://supervisor/core"
export HA_TOKEN="${SUPERVISOR_TOKEN}"

mkdir -p "${COOKIE_DIR}"

echo "[willys] Username: ${WILLYS_USERNAME:0:4}***"
echo "[willys] Sync interval: ${SYNC_INTERVAL}s"

if [ -n "$HA_TOKEN" ]; then
    echo "[willys] SUPERVISOR_TOKEN: ja"
else
    echo "[willys] SUPERVISOR_TOKEN: nej"
fi

exec python3 /app/willys_sync.py
